# etc/redirector.service
[Unit]
Description=redirector

[Service]
ExecStart=/srv/redirector/entrypoint.sh

[Install]
WantedBy=multi-user.target
